# ceph_server/files/etc/cron.monthly/ceph-auto-snapshot
#!/bin/sh

# Only call ceph-auto-snapshot if it's available
which ceph-auto-snapshot > /dev/null || exit 0

exec ceph-auto-snapshot --quiet --syslog --label=monthly --keep=3 home projects
